# supervisor.conf
[program:app]
command=/home/ubuntu/ec2-fastapi/gunicorn_start
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/ubuntu/ec2-fastapi/logs/gunicorn-error.log